# ops/ai-ids-backup.timer
[Unit]
Description=Run IDS backup daily

[Timer]
OnCalendar=*-*-* 03:10:00
Persistent=true

[Install]
WantedBy=timers.target
