# ops/ai-ids-backup.service
[Unit]
Description=AI-IDS DB backup

[Service]
Type=oneshot
WorkingDirectory=%h/AI-IDS/AI-IDS
ExecStart=/usr/bin/bash scripts/backup_db.sh
