# ops/ai-ids-api.service
[Unit]
Description=AI-IDS API (Gunicorn)
After=network.target

[Service]
WorkingDirectory=%h/AI-IDS/AI-IDS
Environment=PYTHONUNBUFFERED=1
Environment=REQUIRE_AUTH=1
ExecStart=/usr/bin/env bash -lc 'source .venv/bin/activate && gunicorn -b 127.0.0.1:5000 -w 2 -k gthread --threads 4 --timeout 60 --access-logfile - api:app'
Restart=on-failure

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectHostname=true
ProtectClock=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

[Install]
WantedBy=default.target
